# crontab
0 12 * * 1 /mnt/data/state-invariants-check/run.bash snapshot_state_summary  
0 12 * * 1 /mnt/data/state-invariants-check/run.bash state_invariants_check
