# crontab
0 12 * * 1 /mnt/data/state-invariants-check/run_task.bash snapshot_state_summary  
0 12 * * 1 /mnt/data/state-invariants-check/run_task.bash state_invariants_check
